# deploy/systemd/lab-inventory-api.service
[Unit]
Description=Lab Inventory API (Fastify + Prisma)
After=network.target docker.service
Requires=network.target

[Service]
Type=simple
User=__APP_USER__
Group=__APP_USER__
Environment=NODE_ENV=development
EnvironmentFile=-/etc/default/lab-inventory
ExecStart=/usr/bin/bash -lc 'set -a; [ -f "${REPO_DIR:-/opt/lab-inventory}/apps/api/.env" ] && . "${REPO_DIR:-/opt/lab-inventory}/apps/api/.env"; set +a; cd "${REPO_DIR:-/opt/lab-inventory}/apps/api" && npm run dev'
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
