# deploy/systemd/lab-inventory-web.service
[Unit]
Description=Lab Inventory Web (Vite Dev Server)
After=network.target lab-inventory-api.service
Requires=network.target

[Service]
Type=simple
User=__APP_USER__
Group=__APP_USER__
WorkingDirectory=/opt/lab-inventory/apps/web
Environment=NODE_ENV=development
Environment=VITE_API_BASE=
ExecStart=/usr/bin/npm run dev -- --host 0.0.0.0 --port 5173
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
